# kvaser_bus_manager/install/kvbm-healthcheck.service.template
[Unit]
Description=KVBM healthcheck (restart on API hang)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Needs permission to restart kvbm.service
User=root
ExecStart={{WORKDIR}}/install/kvbm_healthcheck.sh

# Optional overrides via drop-in or Environment=
Environment=KVBM_HEALTHCHECK_URL=http://127.0.0.1:5000/api/config
Environment=KVBM_HEALTHCHECK_CONNECT_TIMEOUT_S=2
Environment=KVBM_HEALTHCHECK_MAX_TIME_S=3
